# server/systemd/ghostscribe-server.service
[Unit]
Description=GhostScribe STT Server (faster-whisper + FastAPI)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=yojo
Group=yojo
WorkingDirectory=/home/yojo/Documents/code/meridian/ghostscribe/server
EnvironmentFile=/etc/ghostscribe/server.env
ExecStart=/home/yojo/Documents/code/meridian/ghostscribe/server/.venv/bin/uvicorn ghostscribe_server.app:app --host ${GHOSTSCRIBE_HOST} --port ${GHOSTSCRIBE_PORT} --workers 1
Restart=on-failure
RestartSec=3

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=false
ReadWritePaths=/home/yojo/Documents/code/meridian/ghostscribe/server/logs

[Install]
WantedBy=multi-user.target
